# update-motd.d/README
Scripts in this directory are numbered to control ordering.

Any script placed here will be parsed and its output added in order
