# systemd/update-motd.timer
# This file determines how frequently we update the MOTD

[Unit]
Description=Update of the system MOTD

[Timer]
OnBootSec=1min
OnUnitActiveSec=30m
